# my-own-redis.service
[Unit]
Description=My Own redis Service

[Service]
ExecStart=/usr/local/bin/my_own_redis
Restart=always
User=nobody
Group=nogroup
